# scripts/socks-proxy-server.service
[Unit]
Description=WSMQ Service

[Service]
Type=simple
ExecStart=/usr/bin/socks-proxy-server start --port 8080

[Install]
WantedBy=multi-user.target